# ops/cloudflared/cloudflared.service.d/timeout.conf
# systemd drop-in for cloudflared.service.
# Installed by deploy/server.sh ops step to /etc/systemd/system/cloudflared.service.d/.
#
# Why a drop-in (not editing the unit file): cloudflared installs its own
# unit at /etc/systemd/system/cloudflared.service via `cloudflared service
# install`. Overriding the upstream file would lose the override on the
# next cloudflared upgrade. Drop-ins survive.
#
# Why these timeouts:
# - TimeoutStartSec=180 — observed v5.1.0 deploy: tunnel re-establishment
#   takes ~30-90s when the connector reconnects to all four CF edge
#   regions in sequence. The default 90s is enough most of the time but
#   does fail intermittently, which trips `set -e` in deploy/server.sh
#   and aborts the whole release. 180s gives headroom without masking a
#   real hang.
# - TimeoutStopSec=180 — same reasoning: reload sends SIGTERM and waits
#   for graceful drain of in-flight tunneled connections.

[Service]
TimeoutStartSec=180
TimeoutStopSec=180
